# services/user/ibus.service
[Unit]
Description=Ibus keyboard layout

[Service]
ExecStart=/usr/bin/ibus-daemon -xim

[Install]
WantedBy=default.target
